# systemd/user/gcp-snapshot.timer
[Unit]
Description=Ghost Control Plane periodic snapshots

[Timer]
OnBootSec=5min
OnUnitActiveSec=30min
AccuracySec=5min
Persistent=false

[Install]
WantedBy=timers.target
